# microscaler-agent/microscaler-agent.conf
[program:microscaler-agent]
command=/usr/local/microscaler-agent/bin/start_agent.rb

