# liberty/liberty.conf
[program:liberty]
command=/opt/ibm/java/jre/bin/java -Xmx96m -javaagent:/opt/ibm/wlp/bin/tools/ws-javaagent.jar -jar /opt/ibm/wlp/bin/tools/ws-server.jar defaultServer

