# cassandra/cassandra.conf
[program:cassandra]
command=/opt/cassandra/bin/cassandra -f

