# microscaler/controller.conf
[program:controller]
command=/usr/bin/env ruby /usr/local/microscaler/bin/controller.rb
priority=500

