# microscaler/gnatsd.conf
[program:gnatsd]
command=/usr/local/bin/gnatsd
priority=120

