# microscaler/healthmanager.conf
[program:healthmanager]
command=/usr/local/microscaler/bin/hm.rb
priority=500

