# microscaler/mongodb.conf
[program:mongodb]
command=/usr/bin/mongod -f /etc/mongodb.conf
user=mongodb
priority=110

