# microscaler/redis.conf
[program:redis]
command=/usr/bin/redis-server
priority=100
