# microscaler/worker-stop.conf
[program:worker-stop]
command=/usr/bin/env rake --rakefile /usr/local/microscaler/bin/rakefile workers:start QUEUE=stop_queue COUNT=3
priority=500

